# catkinws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/stephan/catkinws/src/ros/rosbash/rostcsh"
