# catkinws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/home/stephan/catkinws/src/ros/rosbash/roszsh"
