# catkinws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/home/stephan/catkinws/src/ros/rosbash/rosbash"
