# catkinws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/home/stephan/catkinws/src/ros/rosbash/rosfish"
