# catkinws/devel/etc/catkin/profile.d/10.rosbuild.sh
/home/stephan/catkinws/devel/.private/rosbuild/etc/catkin/profile.d/10.rosbuild.sh