# catkinws/devel/etc/catkin/profile.d/10.roslaunch.sh
/home/stephan/catkinws/devel/.private/roslaunch/etc/catkin/profile.d/10.roslaunch.sh