# catkinws/devel/etc/catkin/profile.d/15.rosbash.fish
/home/stephan/catkinws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.fish