# catkinws/devel/etc/catkin/profile.d/99.roslisp.sh
/home/stephan/catkinws/devel/.private/roslisp/etc/catkin/profile.d/99.roslisp.sh